# secScanner/secScanner/services/timer_file/secaide.timer
[Unit]
Description=Timer for Secaide.

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d
Unit=secaide.service

[Install]
WantedBy=multi-user.target

