# secScanner/secScanner/services/service_file/secaide.service
[Unit]
Description=AIDE's daemon

[Service]
ExecStart=/usr/sbin/aide --check
Type=simple

